# systemd/nova-hub.service
[Unit]
Description=NOVA Glass Hub (single pane of glass para AI factory)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Ajusta User y WorkingDirectory a tu instalación
User=%i
WorkingDirectory=/home/%i/nova-glass
ExecStart=/usr/bin/python3 -m uvicorn nova_server:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
